# solarprognose_de/rootfs/etc/s6-overlay/s6-rc.d/solarprognose/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Solarprognose.de to MQTT
# Runs solarprognose.sh script
# ==============================================================================

bashio::log.info "Starting Solarprognose..."

exec /usr/bin/solarprognose.sh
